# DI-SNFLK-AzureFunction-MSSQL/kerberos/crontab
SHELL=/bin/bash
*/15 * * * * /usr/local/bin/python /kerberos/login.py
